# ml/launchd/com.nickroot.macroml.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Weekly local ML refresh (Mon 09:00 local; runs on next wake if the
       Mac was asleep). Install:
         mkdir -p ~/Library/LaunchAgents
         cp ml/launchd/com.nickroot.macroml.weekly.plist ~/Library/LaunchAgents/
         launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.nickroot.macroml.weekly.plist
       Remove:
         launchctl bootout gui/$(id -u)/com.nickroot.macroml.weekly -->
  <key>Label</key>
  <string>com.nickroot.macroml.weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nickroot/Developer/macro-tracker/scripts/ml-weekly.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/nickroot/Developer/macro-tracker/ml/out/weekly.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nickroot/Developer/macro-tracker/ml/out/weekly.log</string>
</dict>
</plist>
